# systemd/fireurld.service
[Unit]
Description=fireurl deamon

[Service]
ExecStart=/opt/fireurl/bin/fireurld

[Install]
WantedBy=default.target
